# Pebble_Data_R/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; ip_address:port specifier, *:port for all iface

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; num of main logfile rotation backups; default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200
user=root                 ; default is current user, required if root
environment=
	REDIS_IP="localhost",
	RTGRAPH_IP="localhost",
	RELAY_ID="rtest",
	PEBBLE_DATA_LOC="/media/card/rtest/"

; The rpcinterface:supervisor section must remain in the config file for
; RPC (supervisorctl/web interface) to work.  Additional interfaces may be
; added by defining them in separate [rpcinterface:x] sections.

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; The supervisorctl section configures how supervisorctl will connect to
; supervisord.  configure it match the settings in either the unix_http_server
; or inet_http_server section.

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:pebble_connect_0]
command=python /root/besi-relay-station/pebbledata/pebble_connect.py 0 5000
stopasgroup=true
priority=2
autostart=true
autorestart=true

;[program:pebble_ls_0]
;command=python /root/pebbledata/pebble_log_and_stream.py 0 5000
;priority=3
;autostart=true
;autorestart=true

;[program:pebble_connect_1]
;command=python /root/pebbledata/pebble_connect.py 1 5001
;stopasgroup=true
;priority=3
;autostart=true
;autorestart=true

;[program:pebble_ls_1]
;command=python /root/pebbledata/pebble_log_and_stream.py 1 5001
;priority=2
;autostart=true
;autorestart=true

[program:pebble_bind]
command=/root/besi-relay-station/pebbledata/bind_rfcomm.sh
user=root
priority=1
autostart=true
startsecs=0

;[program:state_monitor]
;command=/root/pebbledata/state_monitor.py
;priority=4
;autostart=true
;autorestart=true

[program:rReadSensors]
command=python /root/besi-relay-station/BESI_LOGGING_R/rReadSensors.py
user=root
priority=3
autostart=true
autorestart=true

[program:rPixie]
command=python /root/besi-relay-station/BESI_LOGGING_R/rPixieLog.py
user=root
priority=1
autostart=true
autorestart=true

[program:rPebbleFeatureExt]
command=python /root/besi-relay-station/BESI_LOGGING_R/MotionFeature.py
user=root
priority=1
autostart=true
autorestart=true

[program:rADC]
command=python /root/besi-relay-station/BESI_LOGGING_R/rawADC.py
user=root
priority=2
autostart=true
autorestart=true

[program:rAudio]
command=python /root/besi-relay-station/BESI_LOGGING_R/Audio.py
user=root
priority=2
autostart=true
autorestart=true
